# roles/rpi_arduino_sensors/templates/fusensors.service.j2
[Unit]
After=tailscaled.service

[Service]
WorkingDirectory={{ rpi_arduino_sensors_git_dir }}
ExecStart={{ rpi_arduino_sensors_git_dir }}/py_venv/bin/python3 -m mqtt_io ./mqtt-io.yml
Restart=always
StandardOutput=append:{{ rpi_arduino_sensors_git_dir }}/poll.log
StandardError=inherit
SyslogIdentifier=fusensors
User=fu
Group=fu

[Install]
WantedBy=multi-user.target
